# src/condor_examples/condor.socket
[Unit]
Description=HTCondor Socket

[Socket]
ListenStream=9618

