# src/condor_scripts/native/DEBIAN/preinst
#!/bin/sh
# preinst script for condor
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


CONDOR_USER=condor
CONDOR_GROUP=condor
GECOS="Condor Daemons"
HOME=/var/lib/condor

stop_service(){
	
	# Stop condor service
	# It will fail if condor is not stopped in time limit
	if [ -x "/etc/init.d/condor" ]; then
	  if which invoke-rc.d >/dev/null 2>&1; then
     		invoke-rc.d condor stop
     	  else
     		/etc/init.d/condor stop
     	  fi
	fi
	
}

create_user(){
	#Add condor group if not existed
	if ! getent passwd $CONDOR_GROUP >/dev/null; then
	  echo "Creating group: $CONDOR_GROUP"
	  addgroup --system $CONDOR_GROUP
	else
	  echo "WARN: Group $CONDOR_GROUP already exists. Will not create it"
	fi

	#Add condor user if not existed
	if ! getent passwd $CONDOR_USER >/dev/null; then
	  echo "Creating user: $CONDOR_USER"
	  #adduser --system --ingroup $CONDOR_GROUP --gecos \"$GECOS\" --no-create-home --home $HOME --disabled-password --disabled-login $CONDOR_USER
	  adduser --system --ingroup $CONDOR_GROUP --no-create-home --home $HOME --disabled-password --disabled-login $CONDOR_USER 
	else
	  echo "WARN: The user $CONDOR_USER already exists. Will not create it"
	fi
}


case "$1" in
    install)
        #stop_service
	create_user
    ;;

    upgrade)
    ;;

    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


