# nmi_tools/mail/fw_branch_status.cron
55 07 * * * /nmi/sbin/cronmail condor-builds@cs.wisc.edu /nmi/bin/fw_branch_status
